# packaging/githome.service
[Unit]
Description=Githome server
Documentation=https://github.com/tamnd/githome
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=githome
Group=githome
EnvironmentFile=/etc/githome/githome.env
ExecStart=/usr/bin/githome
Restart=on-failure
RestartSec=5
WorkingDirectory=/var/lib/githome

# Let the service bind a privileged port if GITHOME_LISTEN_HTTP is set to one,
# without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Sandboxing. The server only needs to write under its data directory.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/githome

[Install]
WantedBy=multi-user.target
